# pacman-db-refresh.service
[Unit]
Description=updating package database

[Service]
Type=simple
ExecStart=/usr/bin/pacman -Sy
